# systemd/46snihdnat-zones.timer
[Unit]
Description=46snihdnat: rebuild zone timer

[Timer]
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
