# systemd/46snihdnat-zones.service
[Unit]
Description=46snihdnat: rebuild zone

[Service]
Type=oneshot
User=snihdnat
ExecStart=/home/snihdnat/bin/46snihdnat-rebuild
